# hermes-paper/systemd/paper-opportunity-monitor.service
[Unit]
Description=TAN paper opportunity monitor (paper-only)
After=paper.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/srv/hermes-os/paper
Environment=PYTHONIOENCODING=utf-8
Environment=PAPER_DATA_DIR=/srv/hermes-os/paper/data
Environment=PAPER_OPPORTUNITY_OUTPUT=/srv/hermes-os/paper/data/paper_opportunity_latest.txt
Environment=PAPER_SEND_DISCORD=1
Environment=PAPER_UV_BIN=/root/.local/bin/uv
ExecStart=/srv/hermes-os/paper/scripts/paper-opportunity-monitor.sh
