# factory/agent-dashboard-cli/yalru-terminal-freeclaude.service
[Unit]
Description=YALRU Paperclip real terminal - FreeClaude
After=network.target yalru-agent-cli.service

[Service]
Type=simple
Restart=always
RestartSec=2
ExecStart=/usr/bin/ttyd -i 127.0.0.1 -p 4293 -b /yalru-terminal/freeclaude -W -t disableLeaveAlert=true -t disableResizeOverlay=true -t titleFixed=FreeClaude /usr/local/bin/yalru-agent-terminal freeclaude

[Install]
WantedBy=multi-user.target
